# roles/configuraciones/files/medicos-voluntarios-deploy-config/deploy/medicos.socket
# Place this file at /etc/systemd/system/medicos.socket
[Unit]
Description=medicos socket

[Socket]
ListenStream=/run/medicos.sock

[Install]
WantedBy=sockets.target
